# cron.d/opennebula-guard
#
# cron.d/opennebula-dnsmasq updates dnsmasq pools from opennebula config
#

*/1 * * * * root /opt/vmstore/scripts/scripts/one-guard.py
